# contrib/wakatime-focusd.service
# wakatime-focusd.service - WakaTime focus tracking daemon
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp wakatime-focusd.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now wakatime-focusd.service
#
# IMPORTANT: Your desktop environment's variables must be visible to systemd
# user services. Most compositors and display managers handle this, but if
# auto-detection fails, add this to your compositor's startup config:
#
#   dbus-update-activation-environment --systemd \
#     WAYLAND_DISPLAY XDG_CURRENT_DESKTOP XDG_RUNTIME_DIR \
#     DISPLAY SWAYSOCK HYPRLAND_INSTANCE_SIGNATURE NIRI_SOCKET
#
# Only the variables relevant to your setup matter.
#
# Logs:
#   journalctl --user -u wakatime-focusd -f

[Unit]
Description=WakaTime focus tracking daemon
Documentation=https://github.com/joshuadavidthomas/wakatime-focusd
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/wakatime-focusd
Restart=on-failure
RestartSec=2

# Optional: set log level
Environment=RUST_LOG=info

# Ensure clean environment
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
